# cron/crontab
*/1 * * * * root curl -f http://localhost:8040/indexator/run
